# ops/backups/systemd/superapp-backup@.timer
[Unit]
Description=Timer for SuperApp DB Backup (%i)

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target

